# unbound/rootfs/etc/s6-overlay/s6-rc.d/init-unbound/run
#!/usr/bin/with-contenv bashio
set -e
bashio::log.info "Starting Unbound DNS add-on..."

# Read options
LISTEN_PORT=$(bashio::app.port '5335/udp')
UNBOUND_VERBOSITY=$(bashio::config 'unbound_verbosity')
DO_IPV6=$(bashio::config 'do_ipv6')
NUM_THREADS=$(bashio::config 'num_threads')
CACHE_MIN_TTL=$(bashio::config 'cache_min_ttl')
CACHE_MAX_TTL=$(bashio::config 'cache_max_ttl')
NO_PERSIST_CACHE=$(bashio::config 'no_persist_cache')

# Cap threads to CPU core count
if bashio::config.has_value 'num_threads'; then
	MAX_THREADS=$(nproc)
	if [ "${NUM_THREADS}" -gt "${MAX_THREADS}" ]; then
		bashio::log.warning "num_threads (${NUM_THREADS}) exceeds available CPU cores (${MAX_THREADS}), capping to ${MAX_THREADS}."
		NUM_THREADS=${MAX_THREADS}
	fi
fi

log_output() {
    bashio::log.info "$(printf '%-20s: %s' "$1" "$2")"
}

log_output "Port" "${LISTEN_PORT}"
log_output "Verbosity" "${UNBOUND_VERBOSITY}"
log_output "IPv6" "${DO_IPV6}"
log_output "Threads" "$(bashio::config.has_value 'num_threads' && echo "${NUM_THREADS}" || echo '1')"
log_output "Min TTL" "$(bashio::config.has_value 'cache_min_ttl' && echo "${CACHE_MIN_TTL}s" || echo '0s')"
log_output "Max TTL" "$(bashio::config.has_value 'cache_max_ttl' && echo "${CACHE_MAX_TTL}s" || echo '86400s')"
log_output "Cache persistence" "$([ "${NO_PERSIST_CACHE}" = "true" ] && echo 'disabled' || echo 'enabled')"

# Persistent config directory
CONFIG_DIR="/config"
UNBOUND_CONF="${CONFIG_DIR}/unbound.conf"
CONF_D_DIR="${CONFIG_DIR}/conf.d"
ZONES_DIR="${CONFIG_DIR}/zones"
ROOT_KEY="/etc/unbound/root.key"
HINTS_FILE="/etc/unbound/root.hints"

mkdir -p "${CONF_D_DIR}" "${ZONES_DIR}"

# Refresh root hints if older than 7 days
if [ ! -f "${HINTS_FILE}" ] || \
	[ "$(find "${HINTS_FILE}" -mtime +7 | wc -l)" -gt 0 ]; then
		bashio::log.info "Refreshing root hints..."
		curl -sSL "https://www.internic.net/domain/named.root" \
		-o "${HINTS_FILE}" \
		|| bashio::log.warning "Could not refresh root hints; using cached copy."
fi

# Detect the current global IPv6 prefix on the LAN interface
GUA_PREFIX=""
NETWORK_INFO=$(curl -s -H "Authorization: Bearer ${SUPERVISOR_TOKEN}" "http://supervisor/network/info" 2>/dev/null || true)
GUA_ADDR=$(echo "${NETWORK_INFO}" | jq -r '.data.interfaces[].ipv6.address[]?' 2>/dev/null | grep -E '^[23][0-9a-f]{3}:' | head -n1 || true)

if [ -n "${GUA_ADDR}" ]; then
    GUA_PREFIX="${GUA_ADDR%%::*}::/64"
    bashio::log.info "Detected global IPv6 prefix: ${GUA_PREFIX}"
else
    bashio::log.warning "Could not detect a global IPv6 prefix via Supervisor API; IPv6 LAN clients may be refused."
fi

OPTIONS_CONF="${CONF_D_DIR}/00-ha-options.conf"
cat > "${OPTIONS_CONF}" << EOF
# Home Assistant GUI config options
server:
    verbosity: ${UNBOUND_VERBOSITY}
	port: ${LISTEN_PORT}
$([ "${UNBOUND_VERBOSITY}" -ge 3 ] && printf '    log-queries: yes\n    log-replies: yes')
$([ "${DO_IPV6}" = "false" ] && echo "    do-ip6: no")
$(bashio::config.has_value 'num_threads' && echo "    num-threads: ${NUM_THREADS}")
$(bashio::config.has_value 'cache_min_ttl' && echo "    cache-min-ttl: ${CACHE_MIN_TTL}")
$(bashio::config.has_value 'cache_max_ttl' && echo "    cache-max-ttl: ${CACHE_MAX_TTL}")
$([ "${DO_IPV6}" = "true" ] && echo "    interface: ::0")
$([ "${DO_IPV6}" = "true" ] && printf '    access-control: ::1/128 allow\n    access-control: fe80::/10 allow\n    access-control: fd00::/8 allow\n    access-control: %s allow\n' "${GUA_PREFIX}")
EOF


# Check for existing unbound.conf and use it if found.
if [ -f "${UNBOUND_CONF}" ]; then
	bashio::log.info "Existing unbound.conf found, validating..."
	if CHECKCONF_OUTPUT=$(unbound-checkconf "${UNBOUND_CONF}" 2>&1); then
		bashio::log.info "Existing configuration validated. Will use that."
		exit 0
	else
		bashio::log.warning "Existing unbound.conf is invalid, regenerating..."
		echo "${CHECKCONF_OUTPUT}"
	fi
fi

# Generate unbound.conf
bashio::log.info "Generating Unbound configuration..."

cat > "${UNBOUND_CONF}" << EOF
server:
	username: ""

	# Interfaces
	interface: 0.0.0.0

	# Access control
	access-control: 0.0.0.0/0 refuse
	access-control: ::0/0 refuse
	access-control: 10.0.0.0/8 allow
	access-control: 127.0.0.0/8 allow
	access-control: 172.16.0.0/12 allow
	access-control: 192.168.0.0/16 allow

	# DNSSEC
	auto-trust-anchor-file: "${ROOT_KEY}"
	root-hints: "${HINTS_FILE}"

	# Privacy & hardening
	hide-identity: yes
	hide-version: yes
	do-not-query-localhost: no
	private-address: 255.255.255.255/32
	private-address: 203.0.113.0/24
	private-address: 198.51.100.0/24
	private-address: 198.18.0.0/15
	private-address: 192.168.0.0/16
	private-address: 192.0.2.0/24
	private-address: 192.0.0.0/24
	private-address: 169.254.0.0/16
	private-address: 172.16.0.0/12
	private-address: 100.64.0.0/10
	private-address: 10.0.0.0/8
	private-address: 0.0.0.0/8
	private-address: fd00::/8
	private-address: fe80::/10
	private-address: 2001:db8::/32
	private-address: ::ffff:0:0/96

	# Caching
	prefetch: yes
	prefetch-key: yes
	neg-cache-size: 4m

	# Performance
	so-rcvbuf: 4m
	rrset-cache-size: 256m
	msg-cache-size: 128m
	serve-expired: yes
	serve-expired-client-timeout: 0
	outgoing-range: 8192
	num-queries-per-thread: 4096

	# Logging
	use-syslog: no

	# Include user-created config files
	include-toplevel: "${CONF_D_DIR}/*.conf"

remote-control:
	control-enable: yes
	control-interface: ${CONFIG_DIR}/unbound.control.sock
	control-use-cert: no

# Hyperlocal
auth-zone:
	name: "."
	master: k.root-servers.net
	master: i.root-servers.net
	fallback-enabled: yes
	for-downstream: no
	for-upstream: yes
	zonefile: "/etc/unbound/root.zone"
EOF

# Validate configuration
bashio::log.info "Validating configuration..."
if CHECKCONF_OUTPUT=$(unbound-checkconf "${UNBOUND_CONF}" 2>&1); then
    bashio::log.info "Configuration validated."
else
    bashio::log.fatal "Unbound configuration is invalid."
    echo "${CHECKCONF_OUTPUT}"
    exit 1
fi

exit 0
