# aiostreams/rootfs/etc/services.d/aiostreams/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Persistent storage
mkdir -p /data
export DATABASE_URI="sqlite:///data/db.sqlite"

# Auto-generate secret key and persist it to /data/secret_key.
bashio::log.info "Checking secret key..."
SECRET_KEY_FILE="/data/secret_key"
if [[ -f "${SECRET_KEY_FILE}" ]]; then
	SECRET_KEY=$(cat "${SECRET_KEY_FILE}")
  bashio::log.info "Using existing auto-generated secret key from ${SECRET_KEY_FILE}."
else
  bashio::log.info "Generating new secret key..."
  SECRET_KEY=$(openssl rand -hex 32)
  echo "${SECRET_KEY}" > "${SECRET_KEY_FILE}"
  chmod 600 "${SECRET_KEY_FILE}"
  bashio::log.info "Generated new secret key and saved to ${SECRET_KEY_FILE}."
fi
export SECRET_KEY="${SECRET_KEY}"

# Map app options to AIOStreams environment variables
bashio::log.info "Reading configuration..."

export LOG_FORMAT=text
export LOG_LEVEL=$(bashio::config 'log_level')
export LOG_TIMEZONE=$(bashio::info.timezone)
export PORT=$(bashio::app.port "3000/tcp")
export ADDON_NAME="AIOStreams"
export ADDON_ID="aiostreams"
export BUILTIN_STREMTHRU_URL="https://stremthru.13377001.xyz"
export LOG_SENSITIVE_INFO="$(bashio::config 'sensitive_info' false)"
export NODE_OPTIONS="--max-old-space-size=500"

# Base URL
bashio::config.has_value 'base_url' && export BASE_URL="$(bashio::config 'base_url')" || export BASE_URL="http://homeassistant.local:${PORT}"

# Dashboard accesss
if bashio::config.has_value 'addon_username' && bashio::config.has_value 'addon_password'; then
	ADDON_USERNAME=$(bashio::config 'addon_username')
	ADDON_PASSWORD=$(bashio::config 'addon_password')
	export AIOSTREAMS_AUTH="${ADDON_USERNAME}:${ADDON_PASSWORD}"
elif bashio::config.has_value 'addon_username' || bashio::config.has_value 'addon_password'; then
	bashio::log.warning "Both username & password must be set to enable dashboard access."
fi

# Trakt Client ID
bashio::config.has_value 'trakt_id' && export TRAKT_CLIENT_ID="$(bashio::config 'trakt_id')" || export FETCH_TRAKT_ALIASES=false

# Anime API
ANIME_API=$(bashio::config 'anime_api')
if [[ "${ANIME_API}" == "false" ]]; then
  export ANIME_DB_LEVEL_OF_DETAIL=none
  export DEFAULT_MAX_CACHE_SIZE=5000
fi

# Trusted UUID
bashio::config.has_value 'trusted_uuid' && export TRUSTED_UUIDS="$(bashio::config 'trusted_uuid')"

# Log basic info
log_output() {
    bashio::log.info "$(printf '%-15s: %s' "$1" "$2")"
}

if [[ "${LOG_LEVEL}" = "error" || "${LOG_LEVEL}" = "warning" ]]; then
	bashio::log.info "Starting AIOStreams..."
	log_output "Base URL" "${BASE_URL}"
	log_output "Port" "${PORT}"
	log_output "Database" "${DATABASE_URI}"
	log_output "Log Level" "${LOG_LEVEL}"
	log_output "Timezone" "${LOG_TIMEZONE}"
	log_output "Anime DB" "${ANIME_API}"
	log_output "Sensitive Info" "${LOG_SENSITIVE_INFO}"
	log_output "Trakt" "$(bashio::config.has_value 'trakt_id' && echo 'configured' || echo 'unconfigured')"
	log_output "Trusted UUID" "$(bashio::config.has_value 'trusted_uuid' && echo 'configured' || echo 'unconfigured')"
fi

# Start AIOStreams
cd /app
exec node packages/server/dist/server.js
