# music_assistant_alexa_beta/rootfs/etc/s6-overlay/s6-rc.d/alexa-skill/run
#!/usr/bin/with-contenv bashio
export TZ=$(bashio::info.timezone)
export AWS_DEFAULT_REGION=$(bashio::config 'AWS_DEFAULT_REGION')
export PORT=$(bashio::config 'PORT')

bashio::log.info "Starting Alexa Skill application using port ${PORT}..."

exec /app/venv/bin/python /app/app/app.py
