# unbound/rootfs/etc/s6-overlay/s6-rc.d/unbound/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting Unbound..."
exec unbound -d -c /config/unbound.conf
